# network-reconnect.timer
[Unit]
Description=Run network-reconnect every five minutes

[Timer]
OnCalendar=*-*-* *:*:0
Persistent=false

[Install]
WantedBy=timers.target
